# moj_linux/etc/systemd/system/log-startup.service
[Unit]
Description=Log startup time

[Service]
Type=oneshot
ExecStart=/home/janezpavel/.local/bin/raba startup

[Install]
WantedBy=multi-user.target
